# addons/home-secrets/rootfs/etc/services.d/home-secrets/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Secrets Server shutdown handler
# ==============================================================================

bashio::log.info "Shutting down Home Secrets Server..."